# brandynette.service
[Unit]
Description=Brandynette Video Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/brandynette
Environment=PORT=6969
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /opt/brandynette/src/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=brandynette

[Install]
WantedBy=multi-user.target
